# ansible/files/flush-tds-logs.service
[Service]
Type=oneshot
ExecStart=/usr/bin/find /var/log/thredds/ -mtime +30 -type f -delete
