# ansible/files/thredds-restart.service
[Service]
Type=oneshot
ExecStart=/metno/thredds-production/thredds-restart.sh
