# ansible/files/thredds-restart.timer.j2
[Timer]
OnCalendar=Mon,Thu *-*-* 08:40:00
[Install]
WantedBy=timers.target
